# templates/etc/systemd/system/traefik.service.j2
[Unit]
Description=traefik
After=network-online.target

[Service]
Type=simple
ExecStart=/opt/traefik/traefik
PIDFile=/run/traefik.pid

[Install]
WantedBy=multi-user.target
